# alice.service
[Unit]
Description=Alice DeMettrie Bot

[Service]
Type=simple
ExecStart=alice

[Install]
WantedBy=multi-user.target
